# qagi-dashboard.service
[Unit]
Description=QAGI Dashboard
After=network.target

[Service]
Type=simple
User=$USER
WorkingDirectory=$(pwd)
ExecStart=/usr/bin/python3 $(pwd)/qagi_dashboard.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
