# qagi.service
[Unit]
Description=QAGI - Quantum GPU-Accelerated General Intelligence
After=network.target

[Service]
Type=simple
User=claudehome
WorkingDirectory=/home/claudehome/claude-workspace/projects/qagi-system
ExecStart=/usr/bin/python3 /home/claudehome/claude-workspace/projects/qagi-system/qagi_core.py
Restart=always
RestartSec=10
StandardOutput=append:/var/log/qagi/qagi.log
StandardError=append:/var/log/qagi/qagi.error.log

[Install]
WantedBy=multi-user.target
